# sshd/sshd.ini
[program:sshd]
priority=10
command=/usr/sbin/sshd -D
autostart=true
startsecs=0
autorestart=true
redirect_stderr=true
stdout_events_enabled=true
stdout_logfile=/var/log/supervisor/supervisord.log
